# meta-amd/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-vr@.service
[Unit]
Description=Flash VR image %I to VR

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/vr-update.sh /tmp/images/%i
