# meta-amd/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-vr-bundle@.service
[Unit]
Description=Flash VR Bundle image %I to VR

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/vr-update /tmp/images/%i
