# meta-amd/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-bp-fw@.service
[Unit]
Description=Flash Backplane FW image %I

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/amd-bmc-bp-update -d /tmp/images/%i
