# meta-amd/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-retimer-bundle@.service
[Unit]
Description=Flash %I Retimer firmware

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/bmc_retimer_update /tmp/images/%i
