# meta-amd/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-retimer@.service
[Unit]
Description=Flash %I Retimer firmware 

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/retimer-update.sh /tmp/images/%i
